# resources/etc/supervisord.d/10_foreman-proxy.conf
[program:foreman-proxy]
command=/start_foreman-proxy.sh
redirect_stderr=true
