# resources/etc/supervisord.d/10_postgresql.conf
[program:postgresql]
command=/start_postgresql.sh
redirect_stderr=true
