# resources/etc/supervisord.d/30_puppet-agent.conf
[program:puppet-agent]
command=/usr/bin/puppet agent --no-daemonize --verbose
redirect_stderr=true
